# arch/00-k81x.rules
ACTION=="add", KERNEL=="hidraw*", SUBSYSTEM=="hidraw", RUN+="/usr/bin/k81x-fkeys -s -u %S%p on"
